# mumo-signal-service.service
[Unit]
Description=Pure Quant Trading Signal Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/evans/mumo-syntax-capital
Environment="PYTHONPATH=/home/evans/mumo-syntax-capital"
ExecStart=/home/evans/mumo-syntax-capital/venv/bin/python /home/evans/mumo-syntax-capital/signal_service.py
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-signals

[Install]
WantedBy=multi-user.target
