# monitoring/mumo-watchdog.service
[Unit]
Description=Signal Service Watchdog (Health Checks)
After=network.target

[Service]
Type=oneshot
User=evans
WorkingDirectory=/home/evans/mumo-syntax-capital
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/evans/mumo-syntax-capital/venv/bin/python /home/evans/mumo-syntax-capital/monitoring/watchdog.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mumo-watchdog
